# make-openwrt/openwrt-files/common-files/etc/init.d/repair_ro
#!/bin/sh /etc/rc.common

START=99

start() {
    bash /sbin/repair_ro >/dev/null 2>&1
}
